# generator/example.crontab
# m h      dom mon dow   command
  0 *      *   *   *     python /home/deploy/public_html/hackernews.overhrd.com/generator/hacker_news_deluxe_rss.py > /home/deploy/public_html/hackernews.overhrd.com/generator/latest.rss ; ruby /home/deploy/public_html/hackernews.overhrd.com/generator/wtf.rb

